# deploy/systemd/trade-cron-jobs.service
[Unit]
Description=Stock Watcher Cron Jobs Service
After=network.target postgresql.service trade-backend.service

[Service]
Type=simple
User=your-user
WorkingDirectory=/path/to/stock-watcher/runtime/cron-jobs
ExecStart=/usr/bin/java -jar -Xms256m -Xmx512m trade-cron-jobs.jar
Environment="SPRING_DATASOURCE_URL=jdbc:postgresql://localhost:5433/trade"
Environment="SPRING_DATASOURCE_USERNAME=postgre"
Environment="SPRING_DATASOURCE_PASSWORD=admin"
Restart=on-failure
RestartSec=10
StandardOutput=append:/path/to/stock-watcher/runtime/logs/cron-jobs.log
StandardError=append:/path/to/stock-watcher/runtime/logs/cron-jobs.error.log

[Install]
WantedBy=multi-user.target

